# deploy/systemd/sparky-reverse-tunnel.service
[Unit]
Description=Reverse SSH tunnel from agent to hub (CI access)
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User=AGENT_USER
ExecStart=/usr/bin/ssh \
    -N \
    -R 2222:localhost:22 \
    -o ServerAliveInterval=60 \
    -o ServerAliveCountMax=3 \
    -o ExitOnForwardFailure=yes \
    -o StrictHostKeyChecking=no \
    -i AGENT_HOME/.ssh/id_ed25519 \
    HUB_SSH_USER@HUB_TAILSCALE_HOST
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
